# device-postfix-el9/src/main/resources/usr/lib/systemd/system/postfix-autodiscovery.target
# systemd service for managing all postfix autoconfiguration.

[Unit]
Description=Postfix Configuration Autodiscovery

[Install]
WantedBy=multi-user.target

