# device-postfix-el9/src/main/resources/etc/systemd/system/postfix.service.d/local.conf
[Unit]
After=postfix-autodiscovery.service
PartOf=mail-autodiscovery.target
Requisite=postfix-autodiscovery.service

[Install]
WantedBy=mail-autodiscovery.target
