# device-postfix-el9/src/main/resources/etc/systemd/system/postfix.service.d/local.conf
[Unit]
After=postfix-autodiscovery.service
Requisite=postfix-autodiscovery.service

